# deploy/drjhagpt.service
# systemd unit for DrJhaGPT Enterprise (Python + venv deployment).
# Edit User / WorkingDirectory / paths, then:
#   sudo cp deploy/drjhagpt.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now drjhagpt

[Unit]
Description=DrJhaGPT Enterprise (Streamlit)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=drjhagpt
WorkingDirectory=/opt/drjhagpt-ent
EnvironmentFile=/opt/drjhagpt-ent/.env
ExecStart=/opt/drjhagpt-ent/.venv/bin/streamlit run streamlit_app.py \
          --server.port=8501 --server.address=0.0.0.0 --server.headless=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
